# deploy/systemd/azcoin-sc-node-payout-scheduler.timer.template
[Unit]
Description=AZCoin SC-node payout scheduler timer (explicit-target manual-approved automation)
Requires=azcoin-sc-node-payout-scheduler.service
# Generated by deploy/scripts/install-azcoin-sc-node-payout-scheduler.sh — do not install raw template.
# Service exit 2 = safe skip (SuccessExitStatus); exit 3 = HALT — not success.

[Timer]
OnCalendar=@AZCOIN_PAYOUT_SCHEDULER_ON_CALENDAR@
Persistent=false
Unit=azcoin-sc-node-payout-scheduler.service

[Install]
WantedBy=timers.target
