# deploy/systemd/azcoin-sc-node-fresh-cycle-automation.timer.template
[Unit]
Description=AZCoin SC-node fresh-cycle payout automation timer
Requires=azcoin-sc-node-fresh-cycle-automation.service
# Generated by deploy/scripts/install-azcoin-sc-node-fresh-cycle-automation.sh

[Timer]
OnCalendar=@AZCOIN_FRESH_CYCLE_AUTOMATION_ON_CALENDAR@
Persistent=false
Unit=azcoin-sc-node-fresh-cycle-automation.service

[Install]
WantedBy=timers.target
