# deploy/systemd/azcoin-sc-node-payout-scheduler.service
[Unit]
Description=AZCoin SC-node payout scheduler (explicit-target manual-approved automation)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=azledger
Group=azledger
WorkingDirectory=/opt/azcoin-super/src/azpool
Environment=PYTHONPATH=/opt/azcoin-super/src/azpool
EnvironmentFile=-/etc/azcoin-super/pool-ledger/collector.env
EnvironmentFile=-/etc/azcoin-super/pool-ledger/payout-scheduler.env
# SAFE BY DEFAULT: no CLI target IDs; scheduler reads env and SAFE_SKIP exit 0 when unset.
# Configure payout-scheduler.env with explicit approved plan/preflight/mode before expecting action.
# Exit codes: 0 = success / SAFE_SKIP; 2 = cadence/gates skip (not a failure);
#   3 = HALT / unsafe — investigate; 1 = malformed config.
SuccessExitStatus=2
ExecStart=/opt/azcoin-super/src/azpool/.venv/bin/python payouts/scripts/sc_node_payout_scheduler.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
